# bot.conf
[program:telegram]
directory=/home/ubuntu/InstaTikTok-Bot
command=python3 main.py

autostart=true
autorestart=true
environment=HOME="/home/ubuntu", USER="ubuntu"

stderr_logfile=/home/ubuntu/InstaTikTok-Bot/logfile_err.log
stdout_logfile=/home/ubuntu/InstaTikTok-Bot/logfile.log
