# bot.conf
[program:telegrambot]
directory=/home/ubuntu/superbot
command=python3 app.py

autostart=true
autorestart=true
environment=HOME="/home/ubuntu", USER="ubuntu"

stderr_logfile=/home/ubuntu/superbot/logfile_err.log
stdout_logfile=/home/ubuntu/superbot/logfile.log
